# openjordy/scripts/podman/openjordy.container.in
# OpenJordy gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into openjordy's ~/.config/containers/systemd/
# {{OPENJORDY_HOME}} is replaced at install time.

[Unit]
Description=OpenJordy gateway (rootless Podman)

[Container]
Image=openjordy:local
ContainerName=openjordy
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{OPENJORDY_HOME}}/.openjordy:/home/node/.openjordy
EnvironmentFile={{OPENJORDY_HOME}}/.openjordy/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
